# scheduler/launchd_tools_list.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd_tools_list.plist — Mac uniquement
     Lance benchmark.js toutes les 2h (heures paires UTC).

     Installation :
       cp scheduler/launchd_tools_list.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/launchd_tools_list.plist

     Désinstallation :
       launchctl unload ~/Library/LaunchAgents/launchd_tools_list.plist
       rm ~/Library/LaunchAgents/launchd_tools_list.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mcp-benchmark.tools-list</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/nico/.nvm/versions/node/v24.14.0/bin/node</string>
    <string>/Users/nico/Projects/remote_mcp_hosting_provider_tools_list_latency_measurement/benchmark.js</string>
    <string>--jitter</string>
    <string>300</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/nico/Projects/remote_mcp_hosting_provider_tools_list_latency_measurement</string>

  <!-- Toutes les 2h aux heures paires UTC -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/nico/Projects/remote_mcp_hosting_provider_tools_list_latency_measurement/logs/tools_list.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/nico/Projects/remote_mcp_hosting_provider_tools_list_latency_measurement/logs/tools_list.log</string>

  <!-- Relance si crash -->
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
